# systemd/wazuh-backup.timer
[Unit]
Description=Weekly Wazuh cold backup

[Timer]
OnCalendar=Sun *-*-* 03:30:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
