# systemd/wz-offsite-pull.timer
[Unit]
Description=Weekly off-box pull of Wazuh backup (after the server-side backup)

[Timer]
OnCalendar=Sun *-*-* 05:30:00
Persistent=true
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
